# systemd/mc-backup@.service
[Unit]
After=mc@%i.service
Description=Back up Minecraft Java Edition server world @ %i
Requisite=mc@%i.service

[Service]
ExecStart=/opt/MC/MCbackup.sh -b /opt/MC/backup_dir /opt/MC/java/%i mc@%i
KillMode=none
Type=oneshot
