# systemd/mc@.socket
[Socket]
ListenFIFO=%t/%N
RemoveOnStop=true
SocketMode=600
SocketUser=mc

[Install]
WantedBy=sockets.target
