# systemd/mc@.socket
[Unit]
BindsTo=%N.service

[Socket]
ListenFIFO=%t/%N
RemoveOnStop=true
SocketMode=600
SocketUser=mc
