# systemd/mcbe-getzip.service
[Unit]
Description=Automatically download Minecraft Bedrock Edition server ZIP

[Service]
ExecStart=/opt/MC/MCBEgetZIP.sh
KillMode=none
Type=oneshot

[Install]
WantedBy=multi-user.target
