# systemd/mcbe-log@.service
[Unit]
After=mcbe@%i.service
Description=Post Minecraft Bedrock Edition server connect/disconnect messages to IRC @ %i
Requisite=mcbe@%i.service

[Service]
ExecStart=/opt/MC/MCBElog.sh -t /tmp/tmux-mc/%i %i
Type=oneshot
User=mc
