# systemd/mcbe-getzip.service
[Unit]
After=network-online.target
Description=Automatically download Minecraft Bedrock Edition server ZIP

[Service]
# DNS check
ExecStartPre=/bin/bash -c 'until host minecraft.net > /dev/null; do sleep 1; done'
# MCBEgetZIP.sh may read y asking if you agree to the EULA
ExecStart=/bin/bash -c 'echo y | /opt/MC/MCBEgetZIP.sh'
KillMode=none
Type=oneshot
User=mc
