# systemd/mcbe-getzip.service
[Unit]
After=network-online.target
Description=Automatically download Minecraft Bedrock Edition server ZIP

[Service]
# DNS check
ExecStartPre=/bin/bash -c 'until ping -c 1 google.com &> /dev/null; do sleep 1; done'
# mcbe_getzip.sh may read y asking if you agree to the EULA
ExecStart=/bin/bash -c 'echo y | /opt/MCscripts/mcbe_getzip.sh'
KillMode=mixed
Type=oneshot
User=mc
