# systemd/mcbe-backup@.service
[Unit]
After=mcbe@%i.service
Before=mcbe-autoupdate@%i.service
Description=Back up Minecraft Bedrock Edition server world @ %i
Requisite=mcbe@%i.service

[Service]
ExecStart=/opt/MC/MCBEbackup.sh -b /opt/MC/backup_dir -t /tmp/tmux-mc/%i /opt/MC/%i %i
KillMode=none
Type=oneshot
User=mc
