# systemd/mcbe-backup@.service
[Unit]
After=mcbe@%i.service
Description=Back up Minecraft Bedrock Edition server world @ %i
Requisite=mcbe@%i.service

[Service]
ExecStart=/opt/MCscripts/bin/mcbe_backup.py /opt/MC/bedrock/%i mcbe@%i
KillMode=mixed
Type=oneshot
