# systemd/mcbe-getzip.timer
[Unit]
Description=Automatically download Minecraft Bedrock Edition server ZIP hourly

[Timer]
OnActiveSec=1 hour

[Install]
WantedBy=timers.target
