# systemd/mcbe-log@.service
[Unit]
After=mcbe@%i.service mcbe-bot@%i.service
Description=Post Minecraft Bedrock Edition server connect/disconnect messages to chat @ %i
PartOf=mcbe@%i.service
Requisite=mcbe@%i.service

[Service]
ExecStart=/opt/MC/MCBElog.sh mcbe@%i
Restart=always
Type=simple

[Install]
# Start after mcbe@%i.service
WantedBy=mcbe@%i.service
