# systemd/mcbe-getzip.timer
[Unit]
Description=Automatically download Minecraft Bedrock Edition server ZIP hourly

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
