# systemd/mcbe-rmbackup@.service
[Unit]
After=mcbe-backup@%i.service
Description=Automatically remove Minecraft Bedrock Edition server world backups more than 2-weeks-old

[Service]
ExecStart=/bin/bash -c '/usr/bin/find "$(readlink -f /opt/MC/backup_dir)/%i_Backups" -type f -mtime +13 -delete'
ExecStartPost=/bin/bash -c '/usr/bin/find "$(readlink -f /opt/MC/backup_dir)/%i_Backups" -type d -empty -delete'
Type=oneshot

[Install]
WantedBy=mcbe-backup.service
# Start after mcbe-backup.service
